# etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

echo 2 > /sys/module/hid_apple/parameters/fnmode

echo
echo "Press any key to clear screen."
read -n1

